# scripts/mac-mini/com.cole.brew-autoupdate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cole.brew-autoupdate</string>

    <!-- Run as cole, since cole owns /opt/homebrew -->
    <key>UserName</key>
    <string>cole</string>
    <key>GroupName</key>
    <string>staff</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/brew-autoupdate.sh</string>
    </array>

    <!-- Sundays at 03:00 local -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/cole/Library/Logs/brew-autoupdate.launchd.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/cole/Library/Logs/brew-autoupdate.launchd.err</string>
</dict>
</plist>
